# xfce/.config/autostart/swapcaps.desktop
[Desktop Entry]
Type=Application
Exec=setxkbmap -option ctrl:swapcaps
Hidden=false
NoDisplay=false
X-GNOME-Autostart-enabled=true
Name=SwapCaps
